# deploy/systemd/secondhand-backend.service.example
[Unit]
Description=Secondhand Platform Django Backend (Daphne ASGI Server)
After=network.target mysql.service redis-server.service
Requires=mysql.service redis-server.service

[Service]
Type=simple
User=你的用户名
Group=www-data
WorkingDirectory=/var/www/secondhand/backend
Environment="PATH=/var/www/secondhand/backend/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=core.settings_production"
ExecStart=/var/www/secondhand/backend/venv/bin/daphne -b 127.0.0.1 -p 8000 core.asgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target




